# contrib/init/onlined.conf
description "Online Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env ONLINED_BIN="/usr/bin/onlined"
env ONLINED_USER="online"
env ONLINED_GROUP="online"
env ONLINED_PIDDIR="/var/run/onlined"
# upstart can't handle variables constructed with other variables
env ONLINED_PIDFILE="/var/run/onlined/onlined.pid"
env ONLINED_CONFIGFILE="/etc/online/online.conf"
env ONLINED_DATADIR="/var/lib/onlined"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # onlined will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$ONLINED_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run onlined."
        echo "The setting must appear in $ONLINED_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Online Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$ONLINED_PIDDIR"
    chmod 0755 "$ONLINED_PIDDIR"
    chown $ONLINED_USER:$ONLINED_GROUP "$ONLINED_PIDDIR"
    chown $ONLINED_USER:$ONLINED_GROUP "$ONLINED_CONFIGFILE"
    chmod 0660 "$ONLINED_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$ONLINED_PIDFILE" \
    --chuid $ONLINED_USER:$ONLINED_GROUP \
    --exec "$ONLINED_BIN" \
    -- \
    -pid="$ONLINED_PIDFILE" \
    -conf="$ONLINED_CONFIGFILE" \
    -datadir="$ONLINED_DATADIR" \
    -disablewallet \
    -daemon

